# overlays/usr/local/etc/rc.d/promscale
. /etc/rc.subr

name="promscale"
rcvar="promscale_enable"

load_rc_config $name

command="/usr/local/bin/promscale"

run_rc_command "$1"