# overlay/usr/local/etc/rc.d/promscale
. /etc/rc.subr

name="promscale"
rcvar=${name}_enable

load_rc_config $name

command=/usr/local/bin/${name}
command_args="-db-name promscale -db-ssl-mode 'allow'"

run_rc_command "$1"