# overlay/usr/local/etc/rc.d/promscale
. /etc/rc.subr

name="promscale"
rcvar="promscale_enable"

load_rc_config $name

command="/usr/local/bin/promscale -db-name promscale -db-ssl-mode 'allow'"

run_rc_command "$1"